# templates/etc/profile.d/script.sh.j2
#!/bin/bash

{{ profile_script.content | default('') }}

# vim: ft=sh:ts=4:sw=4